# .install/dual/autologin@.service
[Service]
ExecStart=-/sbin/agetty -a adam %I 38400
